# trunk/etc/usr/local/etc/periodic/daily/100.qon-cron
#!/bin/sh

# run the Qon cron job

echo "Running Qon cron..."

/www/sites/qon/bin/cron.py

exit 3
