# trunk/etc/usr/local/etc/periodic/daily/800.dump-disks
#!/bin/sh

# determine dump level (see dump(8))
dl=1

# do the dump

echo "Dumping filesystems dump level $dl..."

/sbin/dump -$dl -ua -f /dev/nsa0 /dev/da0s1h

exit 3
