# trunk/etc/www/bin/zeoctl
#!/bin/sh
# ZEO instance control script

# The following two lines are for chkconfig.  On Red Hat Linux (and
# some other systems), you can copy or symlink this script into
# /etc/rc.d/init.d/ and then use chkconfig(8) to automatically start
# ZEO at boot time.

# chkconfig: 345 90 10
# description: start a ZEO server

PYTHON="/www/python/bin/python"
ZEOCTL="/usr/local/bin/zeoctl.py"

CONFIG_FILE="/www/etc/zeo.conf"

exec "$PYTHON" "$ZEOCTL" -C "$CONFIG_FILE" ${1+"$@"}
