# trunk/etc/usr/local/etc/periodic/hourly/200.qon-hourly
#!/bin/sh

# run the Qon cron job

echo "Running Qon hourly cron..."

/www/sites/qon/bin/cron-hourly.py

exit 3
